# deploy/com.jacobmadsen.sgt-surge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jacobmadsen.sgt-surge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jacobmadsen/Projects/Trading/sgt-surge/venv/bin/python</string>
        <string>scripts/run_bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jacobmadsen/Projects/Trading/sgt-surge</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/jacobmadsen/Projects/Trading/sgt-surge</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/jacobmadsen/Projects/Trading/sgt-surge/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jacobmadsen/Projects/Trading/sgt-surge/logs/bot.error.log</string>
</dict>
</plist>
